# Services/FiatStiloSteeringWheelButtons.service
[Unit]
Description=Fiat Stilo Steering Wheel Buttons
Wants=sound.target
After=multi-user.target
[Service]
Type=simple
User=stilo
Group=stilo
WorkingDirectory=/home/stilo/FiatStilo-hudiy/Scripts
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/usr/bin:/bin
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
ExecStart=/home/stilo/FiatStilo-hudiy/Scripts/venv/bin/python3 /home/stilo/FiatStilo-hudiy/Scripts/FiatStiloSteeringWheelButtons.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
[Install]
WantedBy=multi-user.target
